# agenttoyota-ui/agenttoyota.service
# /etc/systemd/system/agenttoyota.service
[Unit]
Description=Toyota React UI (Vite preview on :3000)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/agenttoyota/agenttoyota-ui

# Use Node installed by nvm (adjust version/path if yours differs)
Environment=PATH=/home/ubuntu/.nvm/versions/node/v22.20.0/bin:/usr/local/bin:/usr/bin
Environment=NODE_ENV=production

# Serve the built 'dist' with Vite  on all interfaces, port 
ExecStart=/home/ubuntu/.nvm/versions/node/v22.20.0/bin/npm run dev 

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target